# example/.skel/etc/supervisor/conf.d/fairway.conf
[program:fairway]
directory=/data/fairway
command=/bin/bash run.sh
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
